# amegram.service
[Unit]
After=network.target

[Service]
ExecStart=/home/march/run_amegram.sh

[Install]
WantedBy=default.target
